# bridge/systemd-whatsapp.timer.template
[Unit]
Description=Open Antenna WhatsApp Bridge — poll timer
Requires=open-antenna-whatsapp-bridge.service

[Timer]
OnBootSec=30
OnUnitActiveSec={{POLL_INTERVAL}}s
Unit=open-antenna-whatsapp-bridge.service

[Install]
WantedBy=timers.target
