# bridge/systemd-email.service.template
[Unit]
Description=Open Antenna Email Bridge
After=network.target

[Service]
Type=oneshot
ExecStart={{BRIDGE_SCRIPT}}

Environment=ANTENNA_CONFIG={{CONFIG_PATH}}
Environment=ANTENNA_DATA_DIR={{DATA_DIR}}
Environment=PATH={{PATH_VALUE}}

# stdout/stderr go to journald; also piped to email-bridge.log by the script itself
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
